# Github/res/systemd/gentoo-update-tray.service
[Unit]
Description=Gentoo Update System Tray
Documentation=man:gentoo-update(1)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/bin/gentoo-update --tray
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
